# .config/systemd/user/battery.service
[Unit]
Description=Battery notifications

[Service]
Type=simple
ExecStart=/home/monke/.local/bin/battery
